# data/rpi-monitor.service.in
[Unit]
Description=Raspberry Pi Monitor

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/@PROJECT_NAME@
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
